# supervisor/kafka.conf
[program:kafkaserver]
command=/usr/bin/start-kafka.sh
autostart=true
autorestart=true